# scripts/local-clawgate-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!-- Template: replace YOUR_USERNAME with your macOS account name before loading. -->
<dict>
  <key>Label</key>
  <string>ai.clawgate.watchdog.local</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/YOUR_USERNAME/projects/ios/clawgate/scripts/local-clawgate-watchdog.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/.clawgate/logs/watchdog.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/.clawgate/logs/watchdog.err.log</string>

  <key>ProcessType</key>
  <string>Background</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>WATCHDOG_FUNCTIONAL_CHECK</key>
    <string>0</string>
    <key>WATCHDOG_CONSECUTIVE_FAIL_THRESHOLD</key>
    <string>3</string>
    <key>WATCHDOG_AMBIENT_CHECK</key>
    <string>1</string>
  </dict>
</dict>
</plist>
